# packer/files/wol-enable.rules
# Generated from HashiCorp Packer
ACTION=="add", SUBSYSTEM=="net", NAME=="enp*", RUN+="/usr/sbin/ethtool -s $name wol g"
ACTION=="add", SUBSYSTEM=="net", NAME=="eth*", RUN+="/usr/sbin/ethtool -s $name wol g"
